# ops/ib-gateway/atlas-ib-watchdog.timer
[Unit]
Description=Atlas IB Gateway watchdog (every 10 min)

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
