# services/atlas-telegram-bot.service
[Unit]
Description=Atlas Telegram Approval Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/atlas
ExecStart=/usr/bin/python3 /root/atlas/services/telegram_bot.py
Restart=on-failure
RestartSec=10
Environment=TZ=Australia/Brisbane
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-telegram-bot

[Install]
WantedBy=multi-user.target
