# systemd/atlas-risk-precompute.service
[Unit]
Description=Atlas nightly risk pre-computation (regime matrix, ruin probability, portfolio risk)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
EnvironmentFile=-/etc/atlas/atlas.conf
ExecStart=/usr/bin/python3 /root/atlas/scripts/precompute_risk.py --target=all
StandardOutput=append:/root/atlas/logs/risk_precompute.log
StandardError=append:/root/atlas/logs/risk_precompute.log
TimeoutStartSec=600
