# systemd/atlas-canary-check.timer
[Unit]
Description=Atlas P1.1 Canary Check — daily at 07:00 AEST after all universe sweeps

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true
Unit=atlas-canary-check.service

[Install]
WantedBy=timers.target
