# systemd/atlas-universe-rebuild.timer
[Unit]
Description=Atlas universe rebuild — weekly Sunday 18:00 UTC

[Timer]
# Sunday 18:00 UTC = Monday 04:00 AEST (Brisbane)
# Runs after the week's data is in but before Monday's premarket cron.
OnCalendar=Sun *-*-* 18:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
