# systemd/atlas-research-window@.service
[Unit]
Description=Atlas Research Window — per-universe sweep (%i)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
ExecStart=/root/atlas/scripts/research_window_universe.sh %i
EnvironmentFile=-/etc/atlas/atlas.conf
# Covers worst-case instance (sp500: 4200s sweep + 1500s LLM = 5700s, +3300s buffer = 9000s / 150min).
# Shorter universes exit on completion; wider bound is just a safety net.
# TimeoutAbortSec: if the unit hangs in stop phase, SIGKILL fires after 10min.
TimeoutStartSec=9000
TimeoutStopSec=60
KillMode=control-group
TimeoutAbortSec=600
SendSIGKILL=yes
Nice=10
OOMScoreAdjust=500
EnvironmentFile=-/etc/openclaw/openclaw.env
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=TZ=Australia/Brisbane
