# systemd/atlas-consolidation-closure.service
[Unit]
Description=Atlas consolidation closure -- close sector_etfs+commodity_etfs positions and disable crons (one-shot, 2026-05-05)
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/root/atlas/data/HALT

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
EnvironmentFile=-/etc/atlas/atlas.conf
ExecStart=/root/atlas/scripts/run_consolidation_closure.sh
TimeoutStartSec=1800
StandardOutput=append:/root/atlas/logs/consolidation_closure_2026-05-05.log
StandardError=append:/root/atlas/logs/consolidation_closure_2026-05-05.log

[Install]
WantedBy=multi-user.target
