# systemd/atlas-risk-precompute.timer
[Unit]
Description=Run Atlas risk pre-compute nightly at 22:30 UTC

[Timer]
OnCalendar=*-*-* 22:30:00 UTC
RandomizedDelaySec=300
Persistent=true
Unit=atlas-risk-precompute.service

[Install]
WantedBy=timers.target
