# systemd/atlas-consolidation-closure.timer
[Unit]
Description=Atlas consolidation closure -- fires once at 2026-05-05 14:00 UTC

[Timer]
OnCalendar=2026-05-05 14:00:00 UTC
Persistent=false
RemainAfterElapse=no
Unit=atlas-consolidation-closure.service

[Install]
WantedBy=timers.target
