# systemd/atlas-universe-rebuild.service
[Unit]
Description=Atlas universe rebuild (weekly — all markets)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/atlas
EnvironmentFile=-/etc/atlas/atlas.conf
# Static ETF markets first (fast, no network): commodity_etfs, sector_etfs, etc.
# sp500 last (may take several minutes if yfinance cache is fresh).
# Each ExecStart runs sequentially; if one fails, subsequent ones still run.
ExecStart=/usr/bin/python3 /root/atlas/scripts/rebuild_universe.py --universe commodity_etfs
ExecStart=/usr/bin/python3 /root/atlas/scripts/rebuild_universe.py --universe sector_etfs
ExecStart=/usr/bin/python3 /root/atlas/scripts/rebuild_universe.py --universe defensive_etfs
ExecStart=/usr/bin/python3 /root/atlas/scripts/rebuild_universe.py --universe gold_etfs
ExecStart=/usr/bin/python3 /root/atlas/scripts/rebuild_universe.py --universe treasury_etfs
ExecStart=/usr/bin/python3 /root/atlas/scripts/rebuild_universe.py --universe sp500
StandardOutput=append:/root/atlas/logs/universe_rebuild.log
StandardError=append:/root/atlas/logs/universe_rebuild.log
TimeoutStartSec=1800
# Non-zero exit from any ExecStart fails the unit; log for Telegram health watchdog.
